# deploy/aurum_audit_daily.service
[Unit]
Description=AURUM · daily audit live vs backtest
Documentation=https://github.com/joao/aurum.finance

[Service]
Type=oneshot
User=root
WorkingDirectory=/srv/aurum.finance

# Env file contains AURUM_AUDIT_TG_TOKEN / AURUM_AUDIT_TG_CHAT, sourced
# from the same telegram bot already in config/keys.json. File mode 600.
EnvironmentFile=/etc/aurum/audit_daily.env

ExecStart=/usr/bin/python3 tools/debug/audit_live_vs_backtest_daily.py \
  --days 15 --telegram

# Exit code 2 from the script means "audit ran OK, but there are missed
# signals to investigate". It is a valid successful-completion state;
# without this the service is marked failed every time missed > 0.
SuccessExitStatus=2

StandardOutput=journal
StandardError=journal
SyslogIdentifier=aurum-audit-daily
