# deploy/renaissance_shadow@.service
[Unit]
Description=AURUM · RENAISSANCE shadow runner (%i)
Documentation=https://github.com/joao/aurum.finance
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aurum
WorkingDirectory=/srv/aurum.finance

EnvironmentFile=-/etc/aurum/renaissance-shadow-%i.env
Environment=PYTHONUNBUFFERED=1
Environment=AURUM_RENAISSANCE_SHADOW_LABEL=%i
ExecStart=/usr/bin/python3 tools/maintenance/renaissance_shadow.py --tick-sec 900 --run-hours 0 --label %i

MemoryMax=1G
CPUQuota=150%

Restart=on-failure
RestartSec=30s
TimeoutStopSec=120

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/srv/aurum.finance/data
PrivateTmp=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=aurum-renaissance-shadow-%i

[Install]
WantedBy=multi-user.target
